# haraka-smtp.service
[Unit]
Description=Haraka SMTP
After=network.target

[Service]
ExecStartPre=/bin/mkdir -p /usr/projects/pm2/logs/haraka
ExecStartPre=/bin/chown ubuntu:ubuntu /usr/projects/pm2/logs/haraka
User=root
ExecStart=haraka -c /usr/projects/haraka-smtp
Restart=always
RestartSec=10
StandardOutput=append:/usr/projects/pm2/logs/haraka/app.log
StandardError=append:/usr/projects/pm2/logs/haraka/app-error.log

[Install]
WantedBy=multi-user.target